# scripts/spigot-backup.service
[Unit]
Description=Spigot Backup
After=local-fs.target

[Service]
Type=forking
ExecStart=/usr/local/bin/spigot backup
User=craftbukkit
Group=craftbukkit

[Install]
WantedBy=multi-user.target
